# etc/systemd/system/nvidia-disable.service
[Unit]
Description=Disable NVIDIA card on resume
After=suspend.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'gpu-ctl f'

[Install]
WantedBy=suspend.target