# etc/systemd/system/nvidia-disable.service
[Unit]
Description=Disable NVIDIA card on resume
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'gpu-ctl f'

[Install]
After=systemd-suspend.target