# etc/systemd/system/nvidia-disable.service
[Unit]
Description=Disable NVIDIA card on resume

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'gpu-ctl s && gpu-ctl f && gpu-ctl s'

[Install]
WantedBy=suspend.target sleep.target